# slice files/docker_limit.slice
[Unit]
Description=Slice that limits docker resources
Before=slices.target

[Slice]
#Memory Management
MemoryAccounting=true
MemoryLimit=900M