# docker_settings/docker_limit.slice
[Unit]
Description=Slice that limits docker resources
Before=slices.target

[Slice]
#Memory Management
MemoryAccounting=true
MemoryHigh=600M
MemoryLimit=700M